# systemd-unit-files/rclocal.service
[Unit]
Description=rc.local workaround by maniac
After=multi-user.target

[Service]
ExecStart=/etc/rc.local

[Install]
WantedBy=multi-user.target
